# systemd/tau.timer
[Unit]
Description=Run Tau Time Tracker every minute

[Timer]
OnCalendar=minutely

[Install]
WantedBy=wm.target
